# SystemdServices/ccalpy_gui.service
[Unit]
Description=CCal Gui

[Service]
ExecStart=node /home/__USERNAME__/CCal_V2/WebGUI/server.js
User=root
Restart=on-failure

[Install]
WantedBy=multi-user.target

